# backend/supervisor.conf
[supervisord]
nodaemon=true
user=laravel
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
pidfile = /run/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0760
chown=laravel:laravel

[environment:supervisorctl]
SUPERVISOR_SERVER_URL=unix:///tmp/supervisor.sock

[program:worker]
directory = /var/www/html
process_name = %(program_name)s_%(process_num)02d
command = php artisan queue:work --sleep=3 --tries=3
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true
numprocs = 1
stdout_logfile = /var/log/supervisor/worker.log
redirect_stderr = true
stopwaitsecs = 3600
priority = 6
